# GetCertClient/Resources/Fetch/GoGetCertTask.xml
<?xml version="1.0" encoding="UTF-16"?>
<Task version="1.2" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Author>GoGetCert.com</Author>
    <Description>Executes GoGetCert.exe within a specific security context.</Description>
    <Documentation>
        1) Import "GoGetCertTask.xml" into Task Scheduler.
        2) Set task's security context to a domain admin account.
        3) Add to GoGetCert host: -Task='-OnStartup=True -StartTime="1:10 AM" -CommandEXE=schtasks.exe -CommandArgs="/run /tn GoGetCertTask" -CreateNoWindow'
    </Documentation>
    <URI>\GoGetCertTask</URI>
  </RegistrationInfo>
  <Triggers />
  <Principals>
    <Principal id="Author">
      <UserId>[Domain Admin here.]</UserId>
      <LogonType>Password</LogonType>
      <RunLevel>HighestAvailable</RunLevel>
    </Principal>
  </Principals>
  <Settings>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <AllowHardTerminate>true</AllowHardTerminate>
    <StartWhenAvailable>false</StartWhenAvailable>
    <RunOnlyIfNetworkAvailable>false</RunOnlyIfNetworkAvailable>
    <IdleSettings>
      <StopOnIdleEnd>true</StopOnIdleEnd>
      <RestartOnIdle>false</RestartOnIdle>
    </IdleSettings>
    <AllowStartOnDemand>true</AllowStartOnDemand>
    <Enabled>true</Enabled>
    <Hidden>false</Hidden>
    <RunOnlyIfIdle>false</RunOnlyIfIdle>
    <WakeToRun>false</WakeToRun>
    <ExecutionTimeLimit>PT72H</ExecutionTimeLimit>
    <Priority>7</Priority>
  </Settings>
  <Actions Context="Author">
    <Exec>
      <Command>C:\ProgramData\GoGetCert\GoGetCert.exe</Command>
      <Arguments>-Auto -NoPrompts</Arguments>
      <WorkingDirectory>C:\ProgramData\GoGetCert</WorkingDirectory>
    </Exec>
  </Actions>
</Task>
